# etc/udev/rules.d/99-yubikey-remove.rules
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1050", ATTRS{idProduct}=="0407", ACTION=="remove", RUN+="/home/ahwx/.local/bin/yubikey-state"
