# etc/udev/rules.d/99-yubikey.rules
ACTION=="add", ATTRS{idVendor}=="1050", ATTRS{idProduct}=="0407", RUN+="/home/ahwx/.local/bin/yubikey-state"
